# zfs/cmd/cm/build/buildresult/ceres_cm.service
[Unit]
Description=Start ceres_cm service
DefaultDependencies=no
After=zfs.target zfs-boot.service stmf-svc.service

[Service]
ExecStart=/usr/local/bin/ceres_cm -d 2>/dev/null
Restart=on-failure
StartLimitInterval=1

[Install]
WantedBy=multi-user.target