# zfs/cmd/cm/build/buildresult/ceres_cm.service
[Unit]
Description=Start ceres_cm service
DefaultDependencies=no

[Service]
ExecStart=/usr/local/bin/ceres_cm -d 2>/dev/null
Restart=on-failure

[Install]
WantedBy=multi-user.target