# zfs/cmd/cm/build/buildresult/cm_multi_server.service
[Unit]
Description=Start cm_multi_server service
DefaultDependencies=no

[Service]
ExecStart=/var/cm/script/cm_multi_server.py disable  2>/dev/null
Restart=on-failure

[Install]
WantedBy=multi-user.target