# zfs/etc/systemd/system/configd.service.in
[Unit]
Description=Stmf Configuration Loader Service
DefaultDependencies=no
After=pppt-boot.service
Before=stmf-svc.service

[Service]
ExecStart=@sbindir@/configd
Restart=on-failure

[Install]
WantedBy=zfs.target
