# zfs/etc/systemd/system/fc_boot.service.in
[Unit]
Description=Fibre Channel Kernel Module Loader Service
DefaultDependencies=no
After=stmf-boot.service
After=pppt-boot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe qlf_hengwei

[Install]
WantedBy=zfs.target
