# zfs/etc/systemd/system/fmd.service.in
[Unit]
Description=faulty management service
DefaultDependencies=no
Requires=NetworkManager.service rpcbind.service
After=zfs.target rpcbind.service

[Service]
ExecStart=@sbindir@/fmd
Restart=on-failure

[Install]
WantedBy=multi-user.target
