# zfs/etc/systemd/system/iscsitsvc.service.in
[Unit]
Description=iSCSI Target Activate Service
DefaultDependencies=no
Conflicts=shutdown.target
After=syslog.target
After=stmf-boot.service
After=pppt-boot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/sbin/modprobe ziscsit
ExecStart=@sbindir@/iscsitsvc start
#ExecStop=/path/to/scrip delete-me if not necessary

[Install]
WantedBy=multi-user.target
