# zfs/etc/systemd/system/pppt-boot.service.in
[Unit]
Description=Stmf PPPT Kernel Module Loader Service
DefaultDependencies=no
After=stmf-boot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/modprobe zpppt

[Install]
WantedBy=zfs.target
