# zfs/etc/systemd/system/stmf-svc.service.in
[Unit]
Description=Stmf Management Service
DefaultDependencies=no
After=stmf-boot.service
After=pppt-boot.service
After=configd.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=/bin/sleep 5
ExecStart=@sbindir@/stmfsvc start
ExecStop=@sbindir@/stmfsvc stop

[Install]
WantedBy=multi-user.target
