# zfs/etc/systemd/system/zfs-clusterd.service.in
[Unit]
Description=Start clusterd service
DefaultDependencies=no
Requires=NetworkManager.service zfs-boot.service stmf-svc.service
After=zfs.target zfs-boot.service stmf-svc.service

[Service]
ExecStart=@sbindir@/clusterd -d
Restart=on-failure

[Install]
WantedBy=multi-user.target
